# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

; -----------------------------------------------------------------------
; auth-service — listens on port 8080
; -----------------------------------------------------------------------
[program:auth-service]
command=java -jar /app/auth-service.jar
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/supervisor/auth-service.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/auth-service-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
; Override inter-service URLs so they resolve to localhost inside the container.
; External services (postgres, kafka) are supplied via Railway environment vars.
environment=SERVER_PORT="8080",SPRING_DATASOURCE_URL="%(ENV_SPRING_DATASOURCE_URL)s",SPRING_DATASOURCE_USERNAME="%(ENV_SPRING_DATASOURCE_USERNAME)s",SPRING_DATASOURCE_PASSWORD="%(ENV_SPRING_DATASOURCE_PASSWORD)s",WALLET_SERVICE_URL="http://localhost:8081/api/wallet/createWallet"
priority=100

; -----------------------------------------------------------------------
; wallet-service — listens on port 8081
; auth-service calls this service; within the same container it is
; reachable at localhost:8081.
; -----------------------------------------------------------------------
[program:wallet-service]
command=java -jar /app/wallet-service.jar
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/supervisor/wallet-service.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/wallet-service-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=SERVER_PORT="8081",SPRING_DATASOURCE_URL="%(ENV_SPRING_DATASOURCE_URL)s",SPRING_DATASOURCE_USERNAME="%(ENV_SPRING_DATASOURCE_USERNAME)s",SPRING_DATASOURCE_PASSWORD="%(ENV_SPRING_DATASOURCE_PASSWORD)s"
priority=100

; -----------------------------------------------------------------------
; trade-service — listens on port 8082
; Connects to Kafka and calls wallet-service; within the same container
; wallet-service is reachable at localhost:8081.
; -----------------------------------------------------------------------
[program:trade-service]
command=java -jar /app/trade-service.jar
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/supervisor/trade-service.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/trade-service-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=SERVER_PORT="8082",SPRING_DATASOURCE_URL="%(ENV_SPRING_DATASOURCE_URL)s",SPRING_DATASOURCE_USERNAME="%(ENV_SPRING_DATASOURCE_USERNAME)s",SPRING_DATASOURCE_PASSWORD="%(ENV_SPRING_DATASOURCE_PASSWORD)s",WALLET_SERVICE_URL="http://localhost:8081/api/wallet/"
priority=100

; -----------------------------------------------------------------------
; matching-engine — C++ Kafka consumer / order-matching process
; Reads orders from Kafka and publishes trade results back.
; KAFKA_BROKER is injected from the Railway environment.
; -----------------------------------------------------------------------
[program:matching-engine]
command=/app/MatchingEngine
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/matching-engine.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/matching-engine-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
; KAFKA_BROKER is read directly from the process environment (set by Railway).
; No explicit environment= override needed — supervisord inherits the Railway
; env vars and passes them through automatically.
priority=200

; -----------------------------------------------------------------------
; nginx — serves the React frontend on port 80 and reverse-proxies
; /api/* and /ws to the Java services on localhost.
; -----------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=3
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/nginx-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=50
